# restic-backup.service
[Unit]
Description=Run restic backup locally

[Service]
ExecStart=/home/tanq/bin/restic-backup

[Install]
WantedBy=default.target

